# fn_lock/fn_lock.service
[Unit]
Description=Thinkpad FnLock Fix

[Service]
Type=oneshot
ExecStart=/usr/bin/fn_lock
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
